# examples/systemd/innerwarden-agent.service
[Unit]
Description=Inner Warden - Agent (AI-assisted response)
After=network-online.target innerwarden-sensor.service
Wants=network-online.target
Documentation=https://github.com/InnerWarden/innerwarden

[Service]
Type=simple
User=innerwarden
Group=innerwarden
SupplementaryGroups=systemd-journal
# UMask 0007 keeps new files mode 660 and directories 770. Data in
# /var/lib/innerwarden is security-sensitive (incidents, decisions,
# attacker profiles, dashboard cert/key, graph snapshots). World bits
# must be cleared; group members (sensor runs as root:innerwarden)
# keep read+write for the sensor/agent coordination.
UMask=0007
ExecStart=/usr/local/bin/innerwarden-agent --data-dir /var/lib/innerwarden --config /etc/innerwarden/agent.toml --dashboard --dashboard-bind 0.0.0.0:8787
Restart=on-failure
RestartSec=5

# Environment - API keys, dashboard credentials, Telegram tokens
EnvironmentFile=-/etc/innerwarden/agent.env

# Give it time to flush state and audit trail on shutdown
TimeoutStopSec=15
KillSignal=SIGTERM
SendSIGKILL=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=innerwarden-agent

# SIGHUP triggers config reload (future: hot-reload rules/allowlists)
ExecReload=/bin/kill -HUP $MAINPID

# ── Hardening ──
# NoNewPrivileges omitted: block-ip skills call sudo ufw/iptables/nftables.
# Enable only if using CAP_NET_ADMIN directly or no block-ip skills.
PrivateTmp=yes
ProtectSystem=strict
# /run       - UFW needs /run/ufw.lock
# /etc/ufw   - UFW writes updated rules
ReadWritePaths=/var/lib/innerwarden /run /etc/ufw
ReadOnlyPaths=/etc/innerwarden
# ProtectHome=read-only (NOT =yes): spec-081 managed-agent coexistence must verify
# a co-located AI agent's OWN config in /home (e.g. OpenClaw's /home/lab/.env) by
# stat-ing its owner uid. ProtectHome=yes hides /home entirely, so the verifier
# fails closed and auto-blocks the agent's own activity (observed live on Azure
# 2026-06-21: OpenClaw read its .env then called its own Azure-OpenAI endpoint and
# InnerWarden blocked the endpoint). read-only lets the agent READ /home to verify,
# never write. If you run NO co-located agents, you may set this back to =yes.
ProtectHome=read-only
# Optional, defense-in-depth: CAP_SYS_PTRACE lets the (non-root) agent readlink a
# cross-uid /proc/<pid>/exe for strict interpreter-root verification of a managed
# agent running as another user. Not required — the verifier falls back to argv[0]
# (pinned by the registered cmdline fingerprint) when /proc/exe is unreadable.
# Uncomment to enable strict exe verification:
# AmbientCapabilities=CAP_SYS_PTRACE
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
